# gunicorn.conf
[program:gunicorn]
command=/home/ubuntu/env/bin/gunicorn wsgi
directory=/home/ubuntu
user=nobody
autostart=true
autorestart=true
redirect_stderr=true
