# conf/ubuntu-appserver.conf
[program:gunicorn]
command=gunicorn wsgi
directory=/home/ubuntu
user=nobody
autostart=true
autorestart=true
redirect_stderr=true
